# supervisord.conf
[program:manager]
command= python /home/bae/app/manager.py
numprocess=1
autorestart=true
stdout_logfile=/home/bae/logs/manager.stdout.log
stderr_logfile=/home/bae/logs/manager.stderr.log
